# redbutton.conf
start on usb-device-added ID_VENDOR_ID=1d34 ID_MODEL_ID=000d
stop on usb-device-removed ID_VENDOR_ID=1d34 ID_MODEL_ID=000d

respawn

exec /usr/local/sbin/redbutton
